# motis.service
[Unit]
Description=Motis

[Service]
ExecStart=/home/ubuntu/motis/start.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

